# etc/init/jboss.conf
description "JBoss 7.1.1.Final"

start on (filesystem and networking) or runlevel [2345]
stop on runlevel [016]

env JAVA_HOME=/usr/lib/jvm/java-7-openjdk-amd64
env JBOSS_USER=jboss
env JBOSS_HOME=/opt/jboss
env JBOSS_CONF_F=standalone
env JAVA_OPTS1="-Xms1g -Xmx4g -XX:PermSize=512m -XX:MaxPermSize=1024m -Djava.net.preferIPv4Stack=true -Dorg.jboss.resolver.warning=true"
env JAVA_OPTS2="-Dsun.rmi.dgc.client.gcInterval=3600000 -Dsun.rmi.dgc.server.gcInterval=3600000"
env JAVA_OPTS3="-Djboss.modules.system.pkgs=org.jboss.byteman -Djava.awt.headless=true"
env JAVA_OPTS4="-Djboss.server.default.config=standalone.xml"
env JAVA_OPTS5="-XX:+UseConcMarkSweepGC -XX:+CMSPermGenSweepingEnabled -XX:+CMSClassUnloadingEnabled -XX:+UseCodeCacheFlushing"
env JBOSS_CONFIG=standalone.xml
env JBOSS_SCRIPT=standalone.sh

#console owner
kill timeout 300

# logs to /var/log/upstart/jboss.log
console log
respawn
respawn limit 10 5

setuid jboss
setgid jboss

pre-start script
        echo 'Clear failed deployments mark'
        rm -f ${JBOSS_HOME}/standalone/deployments/*.failed || true

        echo 'Clear vfs'
        rm -f -r ${JBOSS_HOME}/standalone/tmp/vfs/* || true
end script

script
        echo 'Prepare run JBoss'
        # Combine all java_opts into one
        env JAVA_OPTS="$JAVA_OPTS1 $JAVA_OPTS2 $JAVA_OPTS3 $JAVA_OPTS4 $JAVA_OPTS5"
        # Extend paths for other params
        env JBOSS_CONF=${JBOSS_HOME}/bin/${JBOSS_CONF_F}.conf
        # Run JBoss
        echo 'Run JBoss'
        exec ${JBOSS_HOME}/bin/$JBOSS_SCRIPT -c $JBOSS_CONFIG
end script

pre-stop script
        exec ${JBOSS_HOME}/bin/jboss-cli.sh --connect command=:shutdown
        sleep 5s
end script
